# crontab
# Führt den Sync-Job jede Stunde zur Minute 0 aus
# Lädt zuerst die Umgebungsvariablen
# Leitet die System-Logs in eine separate Datei. User-Logs werden intern geschrieben.
0 */8 * * * . /app/cron_env && /usr/local/bin/python /app/sync_all_users.py >> /app/data/system.log 2>&1

# WICHTIG: Stellen Sie sicher, dass diese Datei eine Leerzeile am Ende hat.