# crontab
# Stündlicher Sync (Delta-Sync)
# Supercronic-Format: Umgebungsvariablen werden automatisch vererbt
0 * * * * /usr/local/bin/python /app/sync_all_users.py >> /app/data/system.log 2>&1