# crontab
# Stündlicher Sync (Delta-Sync)
0 * * * * . /app/cron_env && /usr/local/bin/python /app/sync_all_users.py >> /app/data/system.log 2>&1

# Leerzeile am Ende erforderlich